# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/sshd.log
stderr_logfile=/var/log/supervisor/sshd_err.log
priority=10

[program:bore-tunnels]
command=/bin/bash -c "
    while true; do
        bore -s \${BORE_SERVER}:\${BORE_PORT} -p 22 --no-tls-verify &
        bore -s \${BORE_SERVER}:\${BORE_PORT} -p 5000 --no-tls-verify &
        bore -s \${BORE_SERVER}:\${BORE_PORT} -p 5001 --no-tls-verify &
        bore -s \${BORE_SERVER}:\${BORE_PORT} -p 8080 --no-tls-verify &
        bore -s \${BORE_SERVER}:\${BORE_PORT} -p 8888 --no-tls-verify &
        bore -s \${BORE_SERVER}:\${BORE_PORT} -p 9000 --no-tls-verify &
        wait
    done
"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/bore.log
stderr_logfile=/var/log/supervisor/bore_err.log
priority=20
environment=BORE_SERVER="%{ENV_BORE_SERVER}",BORE_PORT="%{ENV_BORE_PORT}"

[program:bot-1]
command=/usr/bin/python3 /opt/bots/examples/bot1.py
directory=/opt/bots/examples
autostart=true
autorestart=true
stdout_logfile=/var/log/bots/bot1.log
stderr_logfile=/var/log/bots/bot1_err.log
priority=30
environment=BOT_TOKEN_1="%{ENV_BOT_TOKEN_1}",BOT_PASSWORD="%{ENV_BOT_PASSWORD}"

[program:bot-2]
command=/usr/bin/python3 /opt/bots/examples/bot2.py
directory=/opt/bots/examples
autostart=true
autorestart=true
stdout_logfile=/var/log/bots/bot2.log
stderr_logfile=/var/log/bots/bot2_err.log
priority=31
environment=BOT_TOKEN_2="%{ENV_BOT_TOKEN_2}"

[program:flask-api]
command=/usr/bin/python3 -m gunicorn --bind 0.0.0.0:8080 --workers 2 app:app
directory=/opt/bots/examples
autostart=true
autorestart=true
stdout_logfile=/var/log/bots/flask.log
stderr_logfile=/var/log/bots/flask_err.log
priority=32