# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/sshd.log
stderr_logfile=/var/log/supervisor/sshd_err.log
priority=10

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx_err.log
priority=15

[program:fail2ban]
command=/usr/bin/fail2ban-server -f
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/fail2ban.log
stderr_logfile=/var/log/supervisor/fail2ban_err.log
priority=16

[program:prometheus]
command=/usr/bin/prometheus --config.file=/etc/prometheus/prometheus.yml --storage.tsdb.path=/var/lib/prometheus --web.enable-lifecycle
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/prometheus.log
stderr_logfile=/var/log/supervisor/prometheus_err.log
priority=17

[program:grafana]
command=/usr/sbin/grafana-server --homepath=/usr/share/grafana cfg:default.paths.logs=/var/log/grafana cfg:default.paths.data=/var/lib/grafana cfg:default.paths.plugins=/var/lib/grafana/plugins
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/grafana.log
stderr_logfile=/var/log/supervisor/grafana_err.log
priority=18

[program:redis]
command=/usr/bin/redis-server --daemonize no
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/redis.log
stderr_logfile=/var/log/supervisor/redis_err.log
priority=19

[program:bore-tunnels]
command=/bin/bash -c "
    while true; do
        bore -s \${BORE_SERVER}:\${BORE_PORT} -p 22 --no-tls-verify &
        bore -s \${BORE_SERVER}:\${BORE_PORT} -p 80 --no-tls-verify &
        bore -s \${BORE_SERVER}:\${BORE_PORT} -p 443 --no-tls-verify &
        bore -s \${BORE_SERVER}:\${BORE_PORT} -p 5000 --no-tls-verify &
        bore -s \${BORE_SERVER}:\${BORE_PORT} -p 5001 --no-tls-verify &
        bore -s \${BORE_SERVER}:\${BORE_PORT} -p 5002 --no-tls-verify &
        bore -s \${BORE_SERVER}:\${BORE_PORT} -p 5003 --no-tls-verify &
        bore -s \${BORE_SERVER}:\${BORE_PORT} -p 8080 --no-tls-verify &
        bore -s \${BORE_SERVER}:\${BORE_PORT} -p 8888 --no-tls-verify &
        bore -s \${BORE_SERVER}:\${BORE_PORT} -p 9000 --no-tls-verify &
        wait
    done
"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/bore.log
stderr_logfile=/var/log/supervisor/bore_err.log
priority=20
environment=BORE_SERVER="%{ENV_BORE_SERVER}",BORE_PORT="%{ENV_BORE_PORT}"

[program:bot-1]
command=/usr/bin/python3 /opt/bots/examples/bot1.py
directory=/opt/bots/examples
autostart=true
autorestart=true
stdout_logfile=/var/log/bots/bot1.log
stderr_logfile=/var/log/bots/bot1_err.log
priority=30
environment=BOT_TOKEN_1="%{ENV_BOT_TOKEN_1}",BOT_PASSWORD="%{ENV_BOT_PASSWORD}",WEBHOOK_URL_1="%{ENV_WEBHOOK_URL_1}"

[program:bot-2]
command=/usr/bin/python3 /opt/bots/examples/bot2.py
directory=/opt/bots/examples
autostart=true
autorestart=true
stdout_logfile=/var/log/bots/bot2.log
stderr_logfile=/var/log/bots/bot2_err.log
priority=31
environment=BOT_TOKEN_2="%{ENV_BOT_TOKEN_2}",BOT_PASSWORD="%{ENV_BOT_PASSWORD}",WEBHOOK_URL_2="%{ENV_WEBHOOK_URL_2}"

[program:bot-3]
command=/usr/bin/python3 /opt/scripts/file-manager/file_manager.py
autostart=true
autorestart=true
stdout_logfile=/var/log/bots/bot3.log
stderr_logfile=/var/log/bots/bot3_err.log
priority=32
environment=BOT_TOKEN_3="%{ENV_BOT_TOKEN_3}",BOT_PASSWORD="%{ENV_BOT_PASSWORD}",WEBHOOK_URL_3="%{ENV_WEBHOOK_URL_3}",FILE_BASE_PATH="/data"

[program:bot-4]
command=/usr/bin/python3 /opt/scripts/system-control/system_control.py
autostart=true
autorestart=true
stdout_logfile=/var/log/bots/bot4.log
stderr_logfile=/var/log/bots/bot4_err.log
priority=33
environment=BOT_TOKEN_4="%{ENV_BOT_TOKEN_4}",BOT_PASSWORD="%{ENV_BOT_PASSWORD}",WEBHOOK_URL_4="%{ENV_WEBHOOK_URL_4}",AUTHORIZED_USERS="%{ENV_AUTHORIZED_USERS}"

[program:webhook-manager]
command=/usr/bin/python3 /opt/scripts/webhook-manager/webhook_manager.py
autostart=true
autorestart=true
stdout_logfile=/var/log/bots/webhook-manager.log
stderr_logfile=/var/log/bots/webhook-manager_err.log
priority=34
environment=BOT_PASSWORD="%{ENV_BOT_PASSWORD}",BOT_TOKEN_1="%{ENV_BOT_TOKEN_1}",BOT_TOKEN_2="%{ENV_BOT_TOKEN_2}",BOT_TOKEN_3="%{ENV_BOT_TOKEN_3}",BOT_TOKEN_4="%{ENV_BOT_TOKEN_4}",WEBHOOK_URL_1="%{ENV_WEBHOOK_URL_1}",WEBHOOK_URL_2="%{ENV_WEBHOOK_URL_2}",WEBHOOK_URL_3="%{ENV_WEBHOOK_URL_3}",WEBHOOK_URL_4="%{ENV_WEBHOOK_URL_4}"

[program:monitoring]
command=python3 /opt/scripts/monitor.py
autostart=true
autorestart=true
stdout_logfile=/var/log/bots/monitor.log
stderr_logfile=/var/log/bots/monitor_err.log
priority=40